# rules/uberAgent-ESA-am-sigma-medium-macos.conf
#
# The rules are generated from the Sigma GitHub repository at https://github.com/SigmaHQ/sigma
# Follow these steps to get the latest rules from the repository with Python
#    1. Clone the repository locally
#    2. Using a commandline, change working directory to the just cloned repository
#    3. Run sigmac -I --target uberagent -r rules/
#
# The rules in this file are marked with sigma-level: medium
#

[ActivityMonitoringRule platform=MacOS]
# Detects when a user manipulates with Firmward Password on MacOS. NOTE - this command has been disabled on silicon-based apple computers.
RuleId = 7ed2c9f7-c59d-4c82-a7e2-f859aa676099
RuleName = Suspicious MacOS Firmware Activity
EventType = Process.Start
Tag = proc-start-suspicious-macos-firmware-activity
RiskScore = 50
Query = (Process.Path == "/usr/sbin/firmwarepasswd" and (Process.CommandLine like r"%setpasswd%" or Process.CommandLine like r"%full%" or Process.CommandLine like r"%delete%" or Process.CommandLine like r"%check%"))

[ActivityMonitoringRule platform=MacOS]
# Detects abuse of the cron utility to perform task scheduling for initial or recurring execution of malicious code. Detection will focus on crontab jobs uploaded from the tmp folder.
RuleId = 7c3b43d8-d794-47d2-800a-d277715aa460
RuleName = Scheduled Cron Task/Job - MacOs
EventType = Process.Start
Tag = proc-start-scheduled-cron-task/job-macos
RiskScore = 50
Annotation = {"mitre_attack": ["T1053.003"]}
Query = (Process.Path like r"%/crontab" and Process.CommandLine like r"%/tmp/%")

[ActivityMonitoringRule platform=MacOS]
# Detect file time attribute change to hide new or changes to existing files
RuleId = 88c0f9d8-30a8-4120-bb6b-ebb54abcf2a0
RuleName = File Time Attribute Change
EventType = Process.Start
Tag = proc-start-file-time-attribute-change
RiskScore = 50
Annotation = {"mitre_attack": ["T1070.006"]}
Query = (Process.Path like r"%/touch" and (Process.CommandLine like r"%-t%" or Process.CommandLine like r"%-acmr%" or Process.CommandLine like r"%-d%" or Process.CommandLine like r"%-r%"))

[ActivityMonitoringRule platform=MacOS]
# Detects deletion of local audit logs
RuleId = acf61bd8-d814-4272-81f0-a7a269aa69aa
RuleName = Indicator Removal on Host - Clear Mac System Logs
EventType = Process.Start
Tag = proc-start-indicator-removal-on-host-clear-mac-system-logs
RiskScore = 50
Annotation = {"mitre_attack": ["T1070.002"]}
Query = ((Process.Path like r"%/rm" or Process.Path like r"%/unlink" or Process.Path like r"%/shred") and (Process.CommandLine like r"%/var/log%" or (Process.CommandLine like r"%/Users/%" and Process.CommandLine like r"%/Library/Logs/%")))

[ActivityMonitoringRule platform=MacOS]
# Detects passwords dumps from Keychain
RuleId = b120b587-a4c2-4b94-875d-99c9807d6955
RuleName = Credentials from Password Stores - Keychain
EventType = Process.Start
Tag = proc-start-credentials-from-password-stores-keychain
RiskScore = 50
Annotation = {"mitre_attack": ["T1555.001"]}
Query = ((Process.Path == "/usr/bin/security" and (Process.CommandLine like r"%find-certificate%" or Process.CommandLine like r"% export %")) or (Process.CommandLine like r"% dump-keychain %" or Process.CommandLine like r"% login-keychain %"))

[ActivityMonitoringRule platform=MacOS]
# Detects usage of system utilities (only grep for now) to discover security software discovery
RuleId = 0ed75b9c-c73b-424d-9e7d-496cd565fbe0
RuleName = Security Software Discovery - MacOs
EventType = Process.Start
Tag = proc-start-security-software-discovery-macos
RiskScore = 50
Annotation = {"mitre_attack": ["T1518.001"]}
Query = (Process.Path == "/usr/bin/grep" and ((Process.CommandLine like r"%nessusd%" or Process.CommandLine like r"%santad%" or Process.CommandLine like r"%CbDefense%" or Process.CommandLine like r"%falcond%" or Process.CommandLine like r"%td-agent%" or Process.CommandLine like r"%packetbeat%" or Process.CommandLine like r"%filebeat%" or Process.CommandLine like r"%auditbeat%" or Process.CommandLine like r"%osqueryd%" or Process.CommandLine like r"%BlockBlock%" or Process.CommandLine like r"%LuLu%") or (Process.CommandLine like r"%Little%" and Process.CommandLine like r"%Snitch%")))

[ActivityMonitoringRule platform=MacOS]
# Detects usage of "find" binary in a suspicious manner to perform discovery
RuleId = 85de3a19-b675-4a51-bfc6-b11a5186c971
RuleName = Potential Discovery Activity Using Find - MacOS
EventType = Process.Start
Tag = proc-start-potential-discovery-activity-using-find-macos
RiskScore = 50
Annotation = {"mitre_attack": ["T1083"]}
Query = (Process.Path like r"%/find" and (Process.CommandLine like r"%-perm -4000%" or Process.CommandLine like r"%-perm -2000%" or Process.CommandLine like r"%-perm 0777%" or Process.CommandLine like r"%-perm -222%" or Process.CommandLine like r"%-perm -o w%" or Process.CommandLine like r"%-perm -o x%" or Process.CommandLine like r"%-perm -u=s%" or Process.CommandLine like r"%-perm -g=s%"))

[ActivityMonitoringRule platform=MacOS]
# Detects when a built-in utility is used to decode and decrypt a payload after a macOS disk image (DMG) is executed. Malware authors may attempt to evade detection and trick users into executing malicious code by encoding and encrypting their payload and placing it in a disk image file. This behavior is consistent with adware or malware families such as Bundlore and Shlayer.
RuleId = 234dc5df-40b5-49d1-bf53-0d44ce778eca
RuleName = Payload Decoded and Decrypted via Built-in Utilities
EventType = Process.Start
Tag = proc-start-payload-decoded-and-decrypted-via-built-in-utilities
RiskScore = 50
Annotation = {"mitre_attack": ["T1059", "T1204", "T1140"]}
Query = (Process.Path like r"%/openssl" and Process.CommandLine like r"%/Volumes/%" and Process.CommandLine like r"%enc%" and Process.CommandLine like r"%-base64%" and Process.CommandLine like r"% -d %")

[ActivityMonitoringRule platform=MacOS]
# Detects commandline operations on shell history files
RuleId = 508a9374-ad52-4789-b568-fc358def2c65
RuleName = Suspicious History File Operations
EventType = Process.Start
Tag = proc-start-suspicious-history-file-operations
RiskScore = 50
Annotation = {"mitre_attack": ["T1552.003"]}
Query = (Process.CommandLine like r"%.bash\_history%" or Process.CommandLine like r"%.zsh\_history%" or Process.CommandLine like r"%.zhistory%" or Process.CommandLine like r"%.history%" or Process.CommandLine like r"%.sh\_history%" or Process.CommandLine like r"%fish\_history%")

[ActivityMonitoringRule platform=MacOS]
# Identifies the execution traces of the XCSSET malware. XCSSET is a macOS trojan that primarily spreads via Xcode projects and maliciously modifies applications. Infected users are also vulnerable to having their credentials, accounts, and other vital data stolen.
RuleId = 47d65ac0-c06f-4ba2-a2e3-d263139d0f51
RuleName = Potential XCSSET Malware Infection
EventType = Process.Start
Tag = proc-start-potential-xcsset-malware-infection
RiskScore = 50
Query = ((Parent.Path like r"%/bash" and Process.Path like r"%/curl" and (Process.CommandLine like r"%/sys/log.php%" or Process.CommandLine like r"%/sys/prepod.php%" or Process.CommandLine like r"%/sys/bin/Pods%") and Process.CommandLine like r"%https://%") or ((Parent.Path like r"%/bash" and Process.Path like r"%/osacompile" and Process.CommandLine like r"%/Users/%" and Process.CommandLine like r"%/Library/Group Containers/%") or (Parent.Path like r"%/bash" and Process.Path like r"%/plutil" and Process.CommandLine like r"%LSUIElement%" and Process.CommandLine like r"%/Users/%" and Process.CommandLine like r"%/Library/Group Containers/%") or (Process.Path like r"%/zip" and Process.CommandLine like r"%-r%" and Process.CommandLine like r"%/Users/%" and Process.CommandLine like r"%/Library/Group Containers/%")))
GenericProperty1 = Parent.Path

[ActivityMonitoringRule platform=MacOS]
# Detects when the macOS Script Editor utility spawns an unusual child process.
RuleId = 6e4dcdd1-e48b-42f7-b2d8-3b413fc58cb4
RuleName = Suspicious Execution via macOS Script Editor
EventType = Process.Start
Tag = proc-start-suspicious-execution-via-macos-script-editor
RiskScore = 50
Annotation = {"mitre_attack": ["T1566", "T1566.002", "T1059", "T1059.002", "T1204", "T1204.001", "T1553"]}
Query = (Parent.Path like r"%/Script Editor" and ((Process.Path like r"%/curl" or Process.Path like r"%/bash" or Process.Path like r"%/sh" or Process.Path like r"%/zsh" or Process.Path like r"%/dash" or Process.Path like r"%/fish" or Process.Path like r"%/osascript" or Process.Path like r"%/mktemp" or Process.Path like r"%/chmod" or Process.Path like r"%/php" or Process.Path like r"%/nohup" or Process.Path like r"%/openssl" or Process.Path like r"%/plutil" or Process.Path like r"%/PlistBuddy" or Process.Path like r"%/xattr" or Process.Path like r"%/sqlite" or Process.Path like r"%/funzip" or Process.Path like r"%/popen") or (Process.Path like r"%python%" or Process.Path like r"%perl%")))
GenericProperty1 = Parent.Path

[ActivityMonitoringRule platform=MacOS]
# Detects execution of AppleScript of the macOS scripting language AppleScript.
RuleId = 1bc2e6c5-0885-472b-bed6-be5ea8eace55
RuleName = MacOS Scripting Interpreter AppleScript
EventType = Process.Start
Tag = proc-start-macos-scripting-interpreter-applescript
RiskScore = 50
Annotation = {"mitre_attack": ["T1059.002"]}
Query = (Process.Path like r"%/osascript" and (Process.CommandLine like r"% -e %" or Process.CommandLine like r"%.scpt%" or Process.CommandLine like r"%.js%"))

[ActivityMonitoringRule platform=MacOS]
# Detects disabling security tools
RuleId = ff39f1a6-84ac-476f-a1af-37fcdf53d7c0
RuleName = Disable Security Tools
EventType = Process.Start
Tag = proc-start-disable-security-tools
RiskScore = 50
Annotation = {"mitre_attack": ["T1562.001"]}
Query = ((Process.Path == "/bin/launchctl" and Process.CommandLine like r"%unload%" and (Process.CommandLine like r"%com.objective-see.lulu.plist%" or Process.CommandLine like r"%com.objective-see.blockblock.plist%" or Process.CommandLine like r"%com.google.santad.plist%" or Process.CommandLine like r"%com.carbonblack.defense.daemon.plist%" or Process.CommandLine like r"%com.carbonblack.daemon.plist%" or Process.CommandLine like r"%at.obdev.littlesnitchd.plist%" or Process.CommandLine like r"%com.tenablesecurity.nessusagent.plist%" or Process.CommandLine like r"%com.opendns.osx.RoamingClientConfigUpdater.plist%" or Process.CommandLine like r"%com.crowdstrike.falcond.plist%" or Process.CommandLine like r"%com.crowdstrike.userdaemon.plist%" or Process.CommandLine like r"%osquery%" or Process.CommandLine like r"%filebeat%" or Process.CommandLine like r"%auditbeat%" or Process.CommandLine like r"%packetbeat%" or Process.CommandLine like r"%td-agent%")) or (Process.Path == "/usr/sbin/spctl" and Process.CommandLine like r"%disable%"))

